# ci/test/magento/etc/supervisor/example.conf
[program:example]
command=python /data/web/supervisor/loop.py
autostart=true
autorestart=true
